# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:docker]
command=/usr/bin/docker daemon
